# .config/runit/sv/syncthing/run
#!/bin/sh
rm log_stdout log_stderr
exec syncthing >> log_stdout 2>> log_stderr
