# .config/runit/sv/seafile-client/finish
#!/bin/sh

seaf-cli stop \
    >> log_stdout 2>> log_stderr
